# root/etc/s6-overlay/s6-rc.d/init-teslacamplayer-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

mkdir -p \
    /media/{SavedClips,SentryClips,RecentClips} &>/dev/null # just so the app starts without errors

# permissions
lsiown -R abc:abc \
    /config \
	/app

# ensure export root is present and writable when explicitly configured
if [ -n "${EXPORT_ROOT_PATH:-}" ]; then
    mkdir -p "$EXPORT_ROOT_PATH" || true
    # align ownership with runtime user (abc respects PUID/PGID from base image)
    lsiown -R abc:abc "$EXPORT_ROOT_PATH"
fi
